# systemd/move-screenshots.service
[Unit]
Description=Move screenshots to strata

[Service]
Type=oneshot
ExecStart=/bin/bash %h/bin/move-screenshots.sh
