# systemd/wallpaper.service
[Unit]
Description=Set random wallpaper

[Service]
Type=oneshot
ExecStart=%h/bin/wallpaper
